# ptxproj/projectroot.cc100/lib/udev/rules.d/97-change-hdmi-settings.rules
KERNEL=="card1", SUBSYSTEM=="drm", ACTION=="change", RUN+="/etc/script/setup_hdmi.sh"

